# files/etc/uci-defaults/99-ap-mode
#!/bin/sh

# Configure the image to boot as a dumb AP on the existing LAN bridge.

LOG_TAG="99-ap-mode"

log_error() {
	local message="$*"

	printf '%s: ERROR: %s\n' "$LOG_TAG" "$message" >&2
	if command -v logger >/dev/null 2>&1; then
		logger -p user.err -t "$LOG_TAG" "$message" 2>/dev/null || true
	fi
}

fail() {
	log_error "$*"
	exit 1
}

uci_set_or_fail() {
	local path="$1"
	local value="$2"

	uci -q set "${path}=${value}" ||
		fail "Failed to set UCI value: ${path}"
}

uci_add_list_or_fail() {
	local path="$1"
	local value="$2"

	uci -q add_list "${path}=${value}" ||
		fail "Failed to add UCI list value: ${path}"
}

uci_delete_if_present() {
	local path="$1"

	if uci -q get "$path" >/dev/null 2>&1; then
		uci -q delete "$path" ||
			fail "Failed to delete UCI value: ${path}"
	fi
}

uci_commit_or_fail() {
	local package="$1"

	uci -q commit "$package" ||
		fail "Failed to commit UCI package: ${package}"
}

configure_wifi_device() {
	local section="$1"
	local band

	WIFI_DEVICE_COUNT=$((WIFI_DEVICE_COUNT + 1))
	config_get band "$section" band

	uci_set_or_fail "wireless.${section}.country" "SE"

	case "$band" in
		2g)
			uci_set_or_fail "wireless.${section}.htmode" "HE40"
			uci_set_or_fail "wireless.${section}.channel" "6"
			;;
		5g)
			uci_set_or_fail "wireless.${section}.htmode" "HE80"
			;;
	esac
}

configure_wifi_iface() {
	local section="$1"
	local mode

	config_get mode "$section" mode
	[ "$mode" = "ap" ] || return 0

	WIFI_AP_COUNT=$((WIFI_AP_COUNT + 1))
	uci_set_or_fail "wireless.${section}.encryption" "psk2"
	uci_set_or_fail "wireless.${section}.key" "password123"
}

main() {
	uci_set_or_fail "network.lan.proto" "static"
	uci_delete_if_present "network.lan.ipaddr"
	uci_add_list_or_fail "network.lan.ipaddr" "10.0.0.9/24"
	uci_delete_if_present "network.lan.netmask"
	uci_delete_if_present "network.lan.ip6assign"
	uci_set_or_fail "network.lan.gateway" "10.0.0.1"
	uci_delete_if_present "network.lan.dns"
	uci_add_list_or_fail "network.lan.dns" "10.0.0.1"
	uci_commit_or_fail "network"

	if ! uci -q get dhcp.lan >/dev/null 2>&1; then
		uci_set_or_fail "dhcp.lan" "dhcp"
		uci_set_or_fail "dhcp.lan.interface" "lan"
	fi

	uci_set_or_fail "dhcp.lan.ignore" "1"
	uci_set_or_fail "dhcp.lan.dhcpv4" "disabled"
	uci_set_or_fail "dhcp.lan.dhcpv6" "disabled"
	uci_set_or_fail "dhcp.lan.ra" "disabled"
	uci_set_or_fail "dhcp.lan.ndp" "disabled"
	uci_commit_or_fail "dhcp"

	# Use the Stockholm timezone, including daylight-saving transitions.
	uci_set_or_fail "system.@system[0].zonename" "Europe/Stockholm"
	uci_set_or_fail "system.@system[0].timezone" \
		"CET-1CEST,M3.5.0,M10.5.0/3"
	uci_commit_or_fail "system"

	# Ensure the hardware-specific MTK wireless sections exist before selecting
	# and customizing them.
	/sbin/wifi config ||
		fail "Wireless hardware discovery failed"

	[ -r /lib/functions.sh ] ||
		fail "Required file is unavailable: /lib/functions.sh"
	. /lib/functions.sh ||
		fail "Failed to load /lib/functions.sh"

	WIFI_DEVICE_COUNT=0
	WIFI_AP_COUNT=0

	config_load wireless ||
		fail "Failed to load the wireless UCI package"
	config_foreach configure_wifi_device wifi-device ||
		fail "Failed while iterating over wireless devices"
	[ "$WIFI_DEVICE_COUNT" -gt 0 ] ||
		fail "No wireless device sections were discovered"

	config_foreach configure_wifi_iface wifi-iface ||
		fail "Failed while iterating over wireless interfaces"
	[ "$WIFI_AP_COUNT" -gt 0 ] ||
		fail "No wireless AP interfaces were discovered"

	uci_commit_or_fail "wireless"
}

main "$@"
