# overlays/ap/files/etc/uci-defaults/99-ap-mode
#!/bin/sh

# Configure the AP image to boot as a dumb AP on the existing LAN bridge.

set -eu

LOG_TAG="99-ap-mode"
COMMON_HELPERS="/lib/custom-builds/uci-defaults-common.sh"

if [ ! -r "$COMMON_HELPERS" ]; then
	printf '%s: ERROR: Required file is unavailable: %s\n' \
		"$LOG_TAG" "$COMMON_HELPERS" >&2
	exit 1
fi

# shellcheck source=/dev/null
. "$COMMON_HELPERS"

main() {
	uci set network.lan.proto='static'
	uci -q delete network.lan.ipaddr || true
	uci add_list network.lan.ipaddr='10.0.0.9/24'
	uci -q delete network.lan.netmask || true
	uci -q delete network.lan.ip6assign || true
	uci set network.lan.gateway='10.0.0.1'
	uci -q delete network.lan.dns || true
	uci add_list network.lan.dns='10.0.0.1'
	uci commit network

	cb_ensure_dhcp_lan_section
	uci set dhcp.lan.ignore='1'
	uci set dhcp.lan.dhcpv4='disabled'
	uci set dhcp.lan.dhcpv6='disabled'
	uci set dhcp.lan.ra='disabled'
	uci set dhcp.lan.ndp='disabled'
	uci commit dhcp

	cb_configure_stockholm_timezone
	cb_apply_wifi_defaults
}

main "$@"
